# Tools/linux/systemd/arducopter.service
[Unit]
Description=ArduCopter
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/arducopter
